# armbian/base/rootfs/etc/systemd/system/grafana-server.service.d/override.conf
[Unit]
After=multi-user.target

[Service]
Restart=always
RestartSec=10
PrivateTmp=true

[Install]
WantedBy=bitboxbase.target
