# armbian/base/rootfs/etc/systemd/system/nginx.service.d/override.conf
[Unit]
After=grafana-server.service startup-checks.service
 
[Service]
Restart=always
RestartSec=10
PrivateTmp=true
