# armbian/base/rootfs/etc/systemd/system/bbbsupervisor.service
[Unit]
Description=BitBox Base Supervisor
After=redis.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/bbbsupervisor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
