# armbian/base/rootfs/etc/systemd/system/bitboxbase.target
[Unit]
Description=BitBox Base startup target
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
