# armbian/base/rootfs/etc/systemd/system/logrotate.timer
[Unit]
Description=BitBox Base logrotation (ever 10 minutes)

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
